# deploy/systemd/bambu-monitor.service
[Unit]
Description=Bambu Monitor backend (FastAPI + MQTT bridge)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bambu
Group=bambu
WorkingDirectory=/opt/bambu-monitor/backend
EnvironmentFile=/etc/bambu-monitor.env
ExecStart=/home/bambu/.local/bin/uv run uvicorn app.main:app --host 127.0.0.1 --port 1111 --workers 1
Restart=on-failure
RestartSec=3s
TimeoutStopSec=10s

# Hardening: the service only needs to read its own dir + write the SQLite DB.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bambu-monitor/backend

[Install]
WantedBy=multi-user.target
