# deploy/systemd/ustreamer.service
[Unit]
Description=ustreamer MJPEG webcam server for Bambu Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bambu
Group=video

# Generic baseline. Per-host customisation (real device path, resolution,
# v4l2-ctl exposure / gain / contrast tweaks) MUST go in a drop-in override:
#   sudo systemctl edit ustreamer.service
# That override lives at /etc/systemd/system/ustreamer.service.d/override.conf
# and is preserved across `update.sh` runs. See docs/DEPLOYMENT.md.
ExecStart=/usr/bin/ustreamer \
  --device=/dev/video0 \
  --resolution=1280x720 \
  --desired-fps=30 \
  --format=MJPEG \
  --host=127.0.0.1 \
  --port=9999 \
  --drop-same-frames=30 \
  --slowdown

Restart=on-failure
RestartSec=3s
TimeoutStopSec=5s

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
